# build_carddb@.service
[Unit]
Description=Update card database (%i)

[Service]
WorkingDirectory=%h/lrrbot
ExecStart=/usr/bin/env pipenv run ./build_carddb_%i.py
