# find_clips_slow.timer
[Unit]
Description=Timer for the find_clips_slow script

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=default.target
