# find_clips.timer
[Unit]
Description=Timer for the find_clips script

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=default.target
