# find_clips.service
[Unit]
Description=Find Twitch clips

[Service]
WorkingDirectory=%h/lrrbot
ExecStart=/usr/bin/env python3 %h/lrrbot/find_clips.py

[Install]
WantedBy=default.target
