# find_clips.service
[Unit]
Description=Find Twitch clips

[Service]
WorkingDirectory=%h/lrrbot
ExecStart=/bin/sh -c "exec %h/lrrbot/venv/bin/python3 find_clips.py day"

[Install]
WantedBy=default.target
